# rootfs/etc/services.d/caddy/run
#!/usr/bin/execlineb -P
with-contenv
s6-setuidgid trojancaddy

/usr/bin/caddy run --config /config/caddy/config.json
